# debian/opsipxeconfd.postinst
#!/bin/sh +e

if [ -x "/etc/init.d/opsipxeconfd" ]; then
	update-rc.d -f opsipxeconfd remove >/dev/null || true
	update-rc.d opsipxeconfd defaults 80 20 >/dev/null
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d opsipxeconfd start || true
	else
		/etc/init.d/opsipxeconfd start || true
	fi
fi

