# debian/opsipxeconfd.postrm
#! /bin/sh -e

if [ "$1" = purge ]; then
	update-rc.d opsipxeconfd remove
fi

#DEBHELPER#

